# packaging/e-voice.service
[Unit]
Description=e-voice LLM dictation processor
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/e-voice daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
